# src/ide-server.service
[Unit]
Description=IDE Server
After=network.target

[Service]
Type=simple
User=droidmaster
WorkingDirectory=/home/droidmaster/ide-server
Environment="DATABASE_URL=sqlite://ide_server.db"
Environment="RUST_LOG=info"
ExecStart=/home/droidmaster/ide-server/target/release/ide-server
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
